# SensorIoT/launchAppScript.conf
env APP_PATH="/home/debian/Sensor-IOT/SensorIoT/comunicationApp"

start on startup  
stop on shutdown

script  
  cd $APP_PATH
  exec pm2 start bin/www
end script 